# ClawMail/systemd/robotmail-heartbeat.service
[Unit]
Description=RobotMail Heartbeat Check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/bin/bash /root/.openclaw/heartbeat-robotmail.sh
